# tools/searxng/searxng.container
[Container]
Image=searxng/searxng:latest
ContainerName=searxng
PublishPort=8080:8080
Volume=%h/searxng-podman/searxng:/etc/searxng:rw,Z
Environment=SEARXNG_BASE_URL=http://localhost:8080/
Environment=SEARXNG_SECRET_KEY={{SECRET_KEY}}
DropCapability=ALL
AddCapability=CHOWN SETGID SETUID
Restart=unless-stopped

[Service]
Restart=always

[Install]
WantedBy=default.target
